# rcouch.conf
[program:rcouch]
command=/rel/rcouch/bin/rcouch start
